# tools/ricoh/daily.crontab
0 9 * * * ccdoperbackup -inc -week CCDOperBackupDisk
